# cmd/opencode.service
[Unit]
Description=OpenCode server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/uve/workspace/opencode
EnvironmentFile=/home/uve/workspace/opencode/.env
ExecStart=/home/uve/.local/bin/opencode serve --port 4096 --hostname 127.0.0.1
Restart=always
RestartSec=3
StartLimitIntervalSec=0

[Install]
WantedBy=default.target
